# 99-usb-drives.rules
# goes in /etc/udev/rules.d/
ACTION=="add", SUBSYSTEM=="block", ENV{DEVTYPE}=="partition", RUN+="/usr/local/bin/mount-usb.sh %k"
